# financial-transcribe.service
[Unit]
Description=Financial Email Transcription Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=codespace
Group=codespace
WorkingDirectory=/workspaces/codespaces-blank
Environment=PATH=/workspaces/codespaces-blank/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /workspaces/codespaces-blank/email_transcribe_financial.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=2G
CPUQuota=80%

[Install]
WantedBy=multi-user.target